# systemd/system/call-history.service
[Unit]
Description=Update Call History

[Service]
Type=oneshot
WorkingDirectory=%h/call-history
ExecStart=poetry run call_history

