# systemd/system/call-history.timer
[Unit]
Description=Update Call History

[Timer]
OnBootSec=1min
OnCalendar=hourly

[Install]
WantedBy=timers.target
